# June/services/june-gpu-multi/supervisord.conf
[supervisord]
nodaemon=true
user=juneuser
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700
chown=juneuser:juneuser

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:june-stt]
command=python /app/stt/main.py
directory=/app/stt
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/stt.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
user=juneuser
environment=PYTHONPATH="/app:/app/stt:/app/tts",PORT="8001",CUDA_VISIBLE_DEVICES="0"
startretries=3
stopsignal=TERM
killasgroup=true

[program:june-tts]
command=python /app/tts/main.py
directory=/app/tts
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/tts.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
user=juneuser
environment=PYTHONPATH="/app:/app/stt:/app/tts",PORT="8000",CUDA_VISIBLE_DEVICES="0"
startretries=3
stopsignal=TERM
killasgroup=true

[group:june-services]
programs=june-stt,june-tts
priority=999