# June/services/june-gpu-multi/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info
minfds=1024
minprocs=200

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:sshd]
command=/usr/sbin/sshd -D
directory=/
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/sshd.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/sshd-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
user=root
startretries=3
startsecs=5
stopsignal=TERM
stopwaitsecs=10
killasgroup=true
priority=10

[program:tailscale]
command=/app/tailscale-connect.sh
directory=/app
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/tailscale.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/tailscale-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
user=juneuser
startretries=5
startsecs=30
stopsignal=TERM
stopwaitsecs=30
killasgroup=true
priority=50

[program:june-stt]
command=python -u /app/stt/main.py
directory=/app/stt
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/stt.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/stt-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
user=juneuser
environment=PYTHONPATH="/app:/app/stt:/app/tts",STT_PORT="8001",CUDA_VISIBLE_DEVICES="0",ALL_PROXY="socks5://localhost:1055/",HTTP_PROXY="http://localhost:1055/",HTTPS_PROXY="http://localhost:1055/",http_proxy="http://localhost:1055/",https_proxy="http://localhost:1055/"
startretries=5
startsecs=45
stopsignal=TERM
stopwaitsecs=30
killasgroup=true
priority=100

[program:june-tts]
command=python -u /app/tts/main.py
directory=/app/tts
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/tts.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/tts-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
user=juneuser
environment=PYTHONPATH="/app:/app/stt:/app/tts",TTS_PORT="8000",CUDA_VISIBLE_DEVICES="0",ALL_PROXY="socks5://localhost:1055/",HTTP_PROXY="http://localhost:1055/",HTTPS_PROXY="http://localhost:1055/",http_proxy="http://localhost:1055/",https_proxy="http://localhost:1055/"
startretries=5
startsecs=60
stopsignal=TERM
stopwaitsecs=30
killasgroup=true
priority=200

[group:june-services]
programs=sshd,tailscale,june-stt,june-tts
priority=999