# June/services/june-vast/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:june-stt]
command=/venv-stt/bin/python /app/stt/main.py
directory=/app/stt
autostart=true
autorestart=true
startretries=3
startsecs=10
stdout_logfile=/app/logs/stt.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/app/logs/stt-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
environment=
    PYTHONUNBUFFERED="1",
    STT_PORT="8001",
    WHISPER_CACHE_DIR="/app/models",
    ORCHESTRATOR_URL="%(ENV_ORCHESTRATOR_URL)s",
    LIVEKIT_WS_URL="%(ENV_LIVEKIT_WS_URL)s",
    LIVEKIT_API_KEY="%(ENV_LIVEKIT_API_KEY)s",
    LIVEKIT_API_SECRET="%(ENV_LIVEKIT_API_SECRET)s",
    CUDA_VISIBLE_DEVICES="0",
    PYTORCH_CUDA_ALLOC_CONF="max_split_size_mb:256"
priority=100
stopwaitsecs=30

[program:june-tts]
command=/venv-tts/bin/python /app/tts/main.py
directory=/app/tts
autostart=true
autorestart=true
startretries=3
startsecs=15
stdout_logfile=/app/logs/tts.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/app/logs/tts-error.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=5
environment=
    PYTHONUNBUFFERED="1",
    TTS_PORT="8000",
    MODEL_DIR="/app/pretrained_models",
    COSYVOICE_MODEL="%(ENV_COSYVOICE_MODEL)s",
    ORCHESTRATOR_URL="%(ENV_ORCHESTRATOR_URL)s",
    LIVEKIT_WS_URL="%(ENV_LIVEKIT_WS_URL)s",
    LIVEKIT_API_KEY="%(ENV_LIVEKIT_API_KEY)s",
    LIVEKIT_API_SECRET="%(ENV_LIVEKIT_API_SECRET)s",
    CUDA_VISIBLE_DEVICES="0",
    PYTORCH_CUDA_ALLOC_CONF="max_split_size_mb:512",
    PYTHONPATH="/opt/CosyVoice/third_party/Matcha-TTS:/opt/CosyVoice:/opt/miniconda/envs/cosyvoice/lib/python3.10/site-packages"
priority=200
stopwaitsecs=30

[group:june-services]
programs=june-stt,june-tts
priority=999