# June/services/june-vast/supervisord copy.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
loglevel=info
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# ==========================================
# Docker Compose Manager
# ==========================================
[program:june-services]
command=docker-compose -f /app/docker-compose.yml up
directory=/app
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/june-services.err.log
stdout_logfile=/var/log/supervisor/june-services.out.log
stderr_logfile_maxbytes=10MB
stdout_logfile_maxbytes=10MB
stderr_logfile_backups=3
stdout_logfile_backups=3
user=root
environment=HOME="/root",USER="root"
stopwaitsecs=30
killasgroup=true
stopasgroup=true

# ==========================================
# Health Check Monitor
# ==========================================
[program:health-monitor]
command=/usr/local/bin/health-monitor.sh
directory=/app
autostart=true
autorestart=true
startretries=999
stderr_logfile=/var/log/supervisor/health-monitor.err.log
stdout_logfile=/var/log/supervisor/health-monitor.out.log
stderr_logfile_maxbytes=5MB
stdout_logfile_maxbytes=5MB
user=root
startsecs=60

# ==========================================
# Log Rotator (Optional)
# ==========================================
[program:log-rotator]
command=/usr/sbin/logrotate -v /etc/logrotate.d/june-services
directory=/app
autostart=true
autorestart=false
startretries=1
stderr_logfile=/var/log/supervisor/log-rotator.err.log
stdout_logfile=/var/log/supervisor/log-rotator.out.log
user=root
startsecs=0